# rootfs/etc/services.d/whatsapp-bridge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Go WhatsApp bridge (whatsmeow client + local REST API).
# Session credentials and the message history SQLite databases are kept on
# the persistent /data volume so re-pairing isn't needed after a restart.
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

mkdir -p /data/store

declare mcp_auth_token=""
if bashio::config.has_value 'mcp_auth_token'; then
    mcp_auth_token=$(bashio::config 'mcp_auth_token')
fi

export WHATSAPP_STORE_DIR="/data/store"
export WHATSAPP_PORT="8080"
export WHATSAPP_QR_PORT="8082"
# Reuse the app's single shared secret to gate the pairing QR image too,
# rather than introducing a second token users have to manage.
export WHATSAPP_QR_TOKEN="${mcp_auth_token}"

cd /app/whatsapp-bridge || bashio::exit.nok "Could not change to /app/whatsapp-bridge"

bashio::log.info "Starting WhatsApp bridge..."
bashio::log.info "If this is the first run, a QR code will appear below - scan it with the WhatsApp app on your phone (Settings > Linked Devices)."
bashio::log.info "If the QR code doesn't render properly in this log viewer: in a browser, go to the SAME address/IP you use for Home Assistant itself, but with port 8082 instead, then add /qr.png to the path. Example: if Home Assistant is at http://192.168.1.50:8123, open http://192.168.1.50:8082/qr.png. This only works while pairing is in progress."
if [ -z "${mcp_auth_token}" ]; then
    bashio::log.warning "mcp_auth_token is not set, so the QR image endpoint above is unauthenticated while pairing - anyone on your network could load it during that window."
else
    bashio::log.info "Since mcp_auth_token is set, add ?token=YOUR_MCP_AUTH_TOKEN to the QR image URL above."
fi

exec ./whatsapp-bridge
